# docs/warranty-agent.service
[Unit]
Description=Warranty Agent System - Multi-agent warranty claims analysis
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/warranty-agent
Environment="PATH=/opt/warranty-agent/.venv/bin"
ExecStart=/opt/warranty-agent/.venv/bin/python -m backend.agent.main server

# Restart configuration
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=warranty-agent

[Install]
WantedBy=multi-user.target
